# deploy/cereblix-stratum.service
[Unit]
Description=Cereblix Stratum bridge (Stratum <-> pool getwork)
After=network-online.target cereblix-pool.service
Wants=network-online.target

[Service]
Type=simple
User=cereblix
Group=cereblix
ExecStart=/opt/cerebra/bin/cereblix-stratum -listen :3333 -pool http://127.0.0.1:18754/api
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
